# packages/KolejkaObserver/debian/python3-kolejkaobserver.service
[Unit]
Description=KolejkaObserver

[Service]
ExecStart=/usr/bin/kolejka-observer

[Install]
WantedBy=multi-user.target
