# kolejka-foreman-system/rc.local
#!/bin/bash
# vim:ts=4:sts=4:sw=4:expandtab

# AUTHORIZED KEYS
mkdir -p /root/.ssh
chmod 755 /root/.ssh
cat > /root/.ssh/authorized_keys <<EOF
ecdsa-sha2-nistp521 AAAAE2VjZHNhLXNoYTItbmlzdHA1MjEAAAAIbmlzdHA1MjEAAACFBAHV7/8MTk8OoNWdXK5Hw7NYMPBoK8XnZVmDnLP7wexBVqqJFenFW75bqqHtgXDmXYVVf7Q9XfVafXM5S1We3YfKbAGYowllBaX3vga8wkvIJj816Yldhv/0Vmy9DXLnzcKAaA90a0LFxZWfsVg6VqMUgDK0ra/CS+ofRzwf07VHKUEr8g== root@random
ecdsa-sha2-nistp521 AAAAE2VjZHNhLXNoYTItbmlzdHA1MjEAAAAIbmlzdHA1MjEAAACFBAB7mwQ+/78MkUKl8sFTINN6Ubyh72iWUeRj6MK82hQuxBp1/QoWlQLlKu+Nw/lcPBB4zGfesD8N+eNHXMx51fdDYAFPfMPUpCM/K5ffzKKF34/OBL5N0QDkhVjDpTNVSiIsDZVGgiKD9kUxr5Qf7IoxXH3KACKFKENSXFavGg90XhYjEQ== gutowski@guliwer
ecdsa-sha2-nistp521 AAAAE2VjZHNhLXNoYTItbmlzdHA1MjEAAAAIbmlzdHA1MjEAAACFBADTp07rSIoLsqtmF1lcIdaP0j8MUjxAeXREQISbNMT/sf88En9ORsgL2q+TL3Uf/isp3yqKnM91Aj4tcECv8IRkZwDnbe74TU/OQQUc63T3cyrJ9Ts0i3DZWLC9NPIVJnafB+QoWVjl65DFb6NthNG13tT+eYkTA3ukCuhUSVnXsILuqA== gutowski@tcs
ecdsa-sha2-nistp521 AAAAE2VjZHNhLXNoYTItbmlzdHA1MjEAAAAIbmlzdHA1MjEAAACFBABE+YFv8KCBjXETHXgd27+2wi2X7x/g5ydmCq36TLPXQVEJZpheBJbfoRZpjiKGmkAp85AuX4QGTOq/mNEhZGI+fAHcLe60D8B4xpXWG/01mdhQa1X2NgDnLMzgduD9vuGWLufbP6M0zC4p27lIT4UceHF/6IYf37WuEbq4MhoB8DnJ0A== root@kolejka
ecdsa-sha2-nistp521 AAAAE2VjZHNhLXNoYTItbmlzdHA1MjEAAAAIbmlzdHA1MjEAAACFBAASuz6dlbQx3L9EVSlyvuq11Lad+vfmUPdY763rSigFBHPocEvVN17M5+ZhxVbBXpp1InQuMFaMkRBn+0VA1jtbeQCeqf3johQctncu+uXCjAmj/Nxv3Q67VIKPvufx5gDpFE/WFpDusc3/aFs0Op6AqjLyHDAGVkM3HK17IXDxIrMhKA== kolejka@matinf.uj.edu.pl
EOF
chmod 644 /root/.ssh/authorized_keys

# SSH HOST KEY
mkdir /etc/ssh
cat > /etc/ssh/ssh_host_ecdsa_key <<EOF
-----BEGIN EC PRIVATE KEY-----
MIHcAgEBBEIAExG5CV2HVHZCBlilqQXFLf5Xr9BKemh/qffULlhMGQ4eyKENFHMt
mQBuj3KzXHHAWRLyJVpH4YV2n01AJp0eJzWgBwYFK4EEACOhgYkDgYYABAAjMGeQ
uShftn/qflWOeTrysOoxCvnSEU1TqlP0BIzHswn9MKHjvUgk+ek7YTiHRqbFzQCn
u0NQbaxQbXn8SLWOBAHmgiutX/KokmFVTsKC4pdtFsk4hAFRoguO0UeAPoXzvoQm
3Aw9ON4EefCVn9s2gdmGfpr4grnEs045IWQjekgprQ==
-----END EC PRIVATE KEY-----
EOF
cat > /etc/ssh/ssh_host_ecdsa_key.pub <<EOF
ecdsa-sha2-nistp521 AAAAE2VjZHNhLXNoYTItbmlzdHA1MjEAAAAIbmlzdHA1MjEAAACFBAAjMGeQuShftn/qflWOeTrysOoxCvnSEU1TqlP0BIzHswn9MKHjvUgk+ek7YTiHRqbFzQCnu0NQbaxQbXn8SLWOBAHmgiutX/KokmFVTsKC4pdtFsk4hAFRoguO0UeAPoXzvoQm3Aw9ON4EefCVn9s2gdmGfpr4grnEs045IWQjekgprQ== kolejka_foreman
EOF
chmod 400 /etc/ssh/ssh_host_ecdsa_key
chmod 644 /etc/ssh/ssh_host_ecdsa_key.pub
sed -e "s|^\( *HostKey .*\)$|\#\1|" -i /etc/ssh/sshd_config
echo "HostKey /etc/ssh/ssh_host_ecdsa_key" >> /etc/ssh/sshd_config
service ssh restart

# SSH SETTINGS FOR GIT
cat > /root/.ssh/config <<EOF
Host bitbucket.org
    StrictHostKeyChecking no
Host github.com
    StrictHostKeyChecking no
EOF

#HTTP PROXY
export HTTP_PROXY="http://149.156.75.219:3128/" FTP_PROXY="http://149.156.75.219:3128/" NO_PROXY="localhost,127.0.0.1,149.156.75,149.156.68,149.156.69"
cat >> /etc/environment <<EOF
HTTP_PROXY="${HTTP_PROXY}"
FTP_PROXY="${FTP_PROXY}"
NO_PROXY="${NO_PROXY}"
EOF

# SUPER USER
groupadd -g 198 -o -f super
useradd -u 198 -g 198 -G sudo -m -o -s /bin/bash -d /home/super -p '$6$7mE8PMaTtClULWgB$uU.l8mt/KNyUHTQJvDavg2j167Y7B9brHf4D7B1xJnA4VXlBspHfCSgOj..j9GeEBo9AMZ2e2P5L8HfN5JJdn/' super

# DOCKER WITH 16 GB SPACE
systemctl stop docker
if [ -e /var/lib/docker/nuke-graph-directory.sh ]; then
    /var/lib/docker/nuke-graph-directory
fi
rm -rf /var/lib/docker
mkdir -p /var/lib/docker
dd if=/dev/zero of=/tmp/docker.xfs bs=$[1024*1024] count=$[16*1024]
mkfs.xfs -f /tmp/docker.xfs
mount -o loop,uquota,pquota /tmp/docker.xfs /var/lib/docker
systemctl start docker

# KOLEJKA
git clone https://github.com/kolejka/kolejka /opt/kolejka
pushd /opt/kolejka
KOLEJKA_SYSTEMS=foreman,observer ./setup.py install
popd

# KOLEJKA OBSERVER DAEMON
systemctl restart kolejka-observer

# KOLEJKA FOREMAN DAEMON
systemctl restart kolejka-foreman
